# conf/novnc.conf
[program:novnc]
command=/root/noVNC/utils/launch.sh --vnc localhost:5901 --listen 6901