# conf/vnc.conf
[program:vnc]
command=/usr/bin/vncserver :1 -geometry 1280x700 -depth 24